# src/main/docker/docker-files/etc/s6/sshd/run
#!/bin/sh
./setup
exec /usr/sbin/sshd -D -e
