# src/main/deb/debian/postrm
#!/bin/sh -e
/usr/sbin/update-rc.d -f git-as-svn remove > /dev/null
