# src/main/deb/debian/postinst
#!/bin/sh -e
/usr/sbin/update-rc.d git-as-svn defaults > /dev/null
/usr/sbin/update-rc.d git-as-svn enable > /dev/null
/etc/init.d/git-as-svn restart
