# deploy/deb/lib/systemd/system/git-as-svn.service
[Unit]
Description=Git as Subversion server
After=network.target syslog.target
ConditionFileNotEmpty=/etc/git-as-svn/git-as-svn.conf

[Service]
User=git
EnvironmentFile=/etc/default/git-as-svn
ExecStart=/usr/bin/git-as-svn --config /etc/git-as-svn/git-as-svn.conf $GAS_ARGS
KillMode=process
Restart=on-failure

[Install]
WantedBy=default.target
