# src/main/docker/docker-files/etc/s6/git-as-svn/run
#!/bin/sh
./setup
cd /data
exec s6-setuidgid $GAS_USER git-as-svn -c /data/config.yml
